# deployments/cloudrun/supervisor-cloudrun/chromium.conf
[program:chromium]
command=/bin/bash -lc '/start-chromium-cloudrun.sh'
user=kernel
autostart=false
autorestart=true
startsecs=5
stdout_logfile=/tmp/supervisord/chromium
redirect_stderr=true